# main/conf/sc_profile_templ.xml
<?xml version="1.0" encoding="US-ASCII"?>
<!DOCTYPE service_bundle SYSTEM "/usr/share/lib/xml/dtd/service_bundle.dtd.1">
<!-- Auto-generated by sysconfig -->
<service_bundle name="sysconfig" type="profile">
  <service name="system/identity" type="service" version="1">
    <instance enabled="true" name="cert"/>
    <instance enabled="true" name="node">
      <property_group name="config" type="application">
        <propval name="nodename" type="astring" value="{node_name}"/>
      </property_group>
    </instance>
  </service>
  <service name="system/name-service/cache" type="service" version="1">
    <instance enabled="true" name="default"/>
  </service>
  <service name="system/name-service/switch" type="service" version="1">
    <property_group name="config" type="application">
      <propval name="default" type="astring" value="files"/>
    </property_group>
    <instance enabled="true" name="default"/>
  </service>
  <service name="system/environment" type="service" version="1">
    <instance enabled="true" name="init">
      <property_group name="environment" type="application">
        <propval name="LANG" type="astring" value="en_US.ISO8859-1"/>
      </property_group>
    </instance>
  </service>
  <service name="system/timezone" type="service" version="1">
    <instance enabled="true" name="default">
      <property_group name="timezone" type="application">
        <propval name="localtime" type="astring" value="US/Eastern"/>
      </property_group>
    </instance>
  </service>
  <service name="system/config-user" type="service" version="1">
    <instance enabled="true" name="default">
      <property_group name="root_account" type="application">
        <propval name="password" type="astring" value="$5$jn8eXfiz$zqxHwPxBc0Uxxxxxxxxxxxxxm7G5ghtFBGoSklg7F8N4"/>
        <propval name="type" type="astring" value="role"/>
        <propval name="login" type="astring" value="root"/>
      </property_group>
      <property_group name="user_account" type="application">
        <propval name="roles" type="astring" value="root"/>
        <propval type="astring" name="shell" value="/usr/bin/bash"/>
        <propval type="astring" name="login" value="admin"/>
        <propval type="astring" name="password" value="$5$6wz1wkGt$RZ.PWs6xBixxxxxxxxxxxxxxf2YuVE2i7HVOrYQib0"/>
        <propval type="astring" name="type" value="normal"/>
        <propval type="astring" name="sudoers" value="ALL=(ALL) ALL"/>
        <propval type="count" name="gid" value="10"/>
        <propval type="astring" name="description" value="Admin"/>
        <propval type="astring" name="profiles" value="System Administrator"/>
      </property_group>
    </instance>
  </service>
  <service name="network/ip-interface-management" type="service" version="1">
    <instance name="default" enabled="true">
      <property_group name="interfaces" type="application">
        <property_group name="lo0" type="interface-loopback">
          <property name="address-family" type="astring">
            <astring_list>
              <value_node value="ipv4"/>
              <value_node value="ipv6"/>
            </astring_list>
          </property>
          <property_group name="v4" type="address-static">
            <propval name="ipv4-address" type="astring" value="127.0.0.1"/>
            <propval name="prefixlen" type="count" value="8"/>
            <propval name="up" type="astring" value="yes"/>
          </property_group>
          <property_group name="v6" type="address-static">
            <propval name="ipv6-address" type="astring" value="::1"/>
            <propval name="prefixlen" type="count" value="128"/>
            <propval name="up" type="astring" value="yes"/>
          </property_group>
        </property_group>
        <property_group name="net0" type="interface-ip">
          <property name="address-family" type="astring">
            <astring_list>
              <value_node value="ipv4"/>
              <value_node value="ipv6"/>
            </astring_list>
          </property>
          <property_group name="v4" type="address-dhcp">
            <propval name="client-id" type="astring" value="{client_id}"/>
            <propval name="dhcp-wait" type="integer" value="-1"/>
            <propval name="primary-interface" type="boolean" value="false"/>
            <propval name="reqhost" type="astring" value=""/>
          </property_group>
          <property_group name="v6" type="address-addrconf">
            <propval name="interface-id" type="astring" value="::"/>
            <propval name="prefixlen" type="count" value="0"/>
            <propval name="stateful" type="astring" value="yes"/>
            <propval name="stateless" type="astring" value="yes"/>
          </property_group>
        </property_group>
      </property_group>
    </instance>
  </service>
</service_bundle>
